# runtimes/common/openrc/runtime.init
#!/sbin/openrc-run
command="/usr/local/bin/runtime"
command_background="no"
description="Run the runtime executable once at boot"